# scripts/yaci-explorer.service
[Unit]
Description=Yaci Explorer Frontend
After=network.target

[Service]
Type=simple
User=yaci
Group=yaci
WorkingDirectory=/opt/yaci-explorer
ExecStart=/usr/local/bin/bun serve ./dist --port 3001
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
